# Webpage Monitor/supervisord.conf
;
;
; General supervisord config
;

[supervisord]
logfile=./supervisord.log ; (main log file;default $CWD/supervisord.log)  ~
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=./supervisord.pid ; (supervisord pidfile;default supervisord.pid)  ~
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

;
;
; Following are necessary for supervisorctl (the interactive shell)
;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

;
;
; Program config
;

[program:page_mon]
command=python2.7 -u ./page_mon.py  
    ; need to run unbuffered (-u) if we want stdout to be flushed
stdout_logfile=./stdout.log
stdout_logfile_maxbytes=1MB
stderr_logfile=./stderr.log
stderr_logfile_maxbytes=1MB
autostart=true
autorestart=true




; To run:
;   set persmissions on page_mon.py
;   supervisord         auto-finds supervisord.conf and runs it
;   supervisorctl       interactive prompt
;       stop page_mon
;       restart page_mon
;       shutdown page_mon
;       help
;   supervisorctl shutdown     shutdown everything
; More info:
;   http://bixly.com/blog/supervisord-or-how-i-learned-to-stop-worrying-and-um-use-supervisord/
;   http://supervisord.org/running.html
;   supervisord -n          debug mode